# crontab
59 23 */2 * * /usr/bin/php artisan transactions:summ >/dev/null 2>&1